# mo-security-mcp.service
[Unit]
Description=WordPress Code Review MCP Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node build/index.js
Restart=always

[Install]
WantedBy=multi-user.target 